# configs/cart.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.learndevopseasy.online
Environment=CATALOGUE_HOST=catalogue.learndevopseasy.online
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target